# infra/systemd/gunicorn.service
[Unit]
Description=Elbrus Cloud Gunicorn daemon
Requires=gunicorn.socket
After=network.target

[Service]
Type=notify
User=elbrus
Group=elbrus
RuntimeDirectory=elbrus
RuntimeDirectoryPreserve=yes
WorkingDirectory=/opt/elbrus/app/backend
EnvironmentFile=/opt/elbrus/app/.env
Environment=DJANGO_SETTINGS_MODULE=elbrus.settings.prod
ExecStart=/opt/elbrus/venv/bin/gunicorn \
          --workers 3 \
          --bind unix:/run/elbrus/gunicorn.sock \
          --access-logfile - \
          --error-logfile - \
          --timeout 60 \
          --limit-request-line 8190 \
          --limit-request-fields 100 \
          --limit-request-field_size 8190 \
          --max-requests 1000 \
          --max-requests-jitter 100 \
          elbrus.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
Restart=on-failure
RestartSec=5

# --- Sandbox / hardening ----------------------------------------------------
# Namespace-based directives (PrivateTmp, PrivateDevices, ProtectSystem,
# ProtectHome, Protect{Kernel,Control,Clock,Hostname}*, RestrictNamespaces)
# are intentionally omitted: they require mount/UTS namespace kernel features
# that Azure B-series VMs restrict for non-root services, causing
# status=226/NAMESPACE at startup. The directives below provide meaningful
# hardening without namespace support.
NoNewPrivileges=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @debug @cpu-emulation @obsolete
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0027

[Install]
WantedBy=multi-user.target
